# .docker/php/supervisor.conf
[supervisord]
nodaemon=true

[program:queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --queue=default,notifications --sleep=3 --tries=3
autostart=true
autorestart=true
numprocs=1
user=developer
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/queue-worker.log